# nginx/crontab
0 0 * * * /bin/cron-acme.sh >> /var/log/acme.log
